# start-uptime.service
[Unit]
Description=my startup script
[Service]
ExecStart=/home/mohamedkhedr/Bash_script/file_services.sh
[Install]
WantedBy=multi-user.target